# templates/default/faye.service.erb
[Unit]
Description=Faye

[Service]
type=forking
User=<%= node['faye']['user'] %>
EnvironmentFile=<%= "#{node['faye']['systemd']['env_path']}/faye" %>
PIDFile=/var/run/<%= node['faye']['app_name'] %>.pid
ExecStart=/usr/bin/node <%= node['faye']['path'] %>/current/server.js

[Install]
WantedBy=multi-user.target
